# supervisor-app.conf
[program:django-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/code/uwsgi.ini

[program:nginx-django]
command = /usr/sbin/nginx
